# sogo-update-autoreply.timer
[Unit]
Description=Automatically update SOGo auto replies

[Timer]
OnUnitActiveSec=5m

[Install]
WantedBy=timers.target
