# sogo-expire-sessions.timer
[Unit]
Description=Automatically clean up expired SOGO sessions

[Timer]
OnUnitActiveSec=1m

[Install]
WantedBy=timers.target
