# sogo-email-alarms.timer
[Unit]
Description=Automatically send SOGO email alarms

[Timer]
OnUnitActiveSec=1m

[Install]
WantedBy=timers.target
