# sogo-email-alarms.service
[Unit]
Description=Send SOGO email alarms

[Service]
Type=oneshot
ExecStart=/usr/sbin/sogo-ealarms-notify
User=sogo
