# sogo-update-autoreply.service
[Unit]
Description=Update SOGo autoreplies

[Service]
Type=oneshot
ExecStart=/usr/bin/sogo-tool update-autoreply -p /etc/sogo/sieve.creds
User=sogo
